# dockerfiles/dockerfiles_win/supervisor_win.conf
[supervisord]
nodaemon=true

#[program:rsyslog]
#command=/bin/bash -c "/usr/sbin/rsyslogd -n"

#[program:iwatch]
#command:/bin/bash -c "iwatch -e close,create,delete,delete_self,modify,move,moved_from,moved_to,open -s -r /home /media"
#command:/bin/bash -c "iwatch -e default,open,modify -s -r /home /media"

#[program:filebeat]
#command:/bin/bash -c "chmod go-w /etc/filebeat.yml && /usr/share/filebeat/bin/filebeat -e -c /etc/filebeat/filebeat.yml"

[program:auditbeat]
command:/bin/bash -c "chmod go-w /etc/auditbeat/auditbeat_ubuntu.yml && auditbeat -e -c /etc/auditbeat/auditbeat_ubuntu.yml"

[program:sshd]
command=/bin/bash -c "/usr/sbin/sshd -D"