# dockerfiles/dockerfiles_auditbeat/supervisor.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/bin/bash -c "/usr/sbin/sshd -D"

[program:rsyslog]
command=/bin/bash -c "/usr/sbin/rsyslogd -n"

[program:iwatch]
#command:/bin/bash -c "iwatch -e close,create,delete,delete_self,modify,move,moved_from,moved_to,open -s -r /home /media"
command:/bin/bash -c "iwatch -e open -s -r /home"

[program:auditbeat]
command:/bin/bash -c "auditbeat -e -c /etc/auditbeat/auditbeat.yml"

[program:filebeat]
command:/bin/bash -c "/opt/filebeat/filebeat -e -c /etc/filebeat/filebeat.yml && tail -f /dev/null"

#[program:auditd]
#command=/bin/bash -c "/etc/init.d/auditd start && tail -f /dev/null"

#[program:dropbox]
#command=/bin/bash -c "/root/run"
